# .config/systemd/user/mu.timer
[Unit]
Description=Timer for mu indexing

[Timer]
OnUnitInactiveSec=15m

[Install]
WantedBy=mail.target
